# ocsrv-supervisor.ini
[program:ocserv]
command=/usr/local/sbin/ocserv -c /config/ocserv.conf -f
autostart=true
autorestart=true
stderr_logfile=/var/log/ocserv.err.log
stdout_logfile=/var/log/ocserv.log
environment=LISTEN_PORT=443